# DEBIAN/postinst
#!/bin/sh

systemctl enable /etc/systemd/system/gauge-control.service
systemctl restart gauge-control.service
